# connman-adapter.upstart
# -*- mode: conf; -*-

start on stopped finish

respawn
exec /usr/bin/connman-adapter

